# deploy/lugrouplib.service
[Unit]
Description=LuGroupLib KiCad library web app
After=network.target

[Service]
Type=simple
# Adjust User and paths to match your VPS.
User=lugrouplib
WorkingDirectory=/opt/lugrouplib/server
Environment=PORT=8000
Environment=LUGROUPLIB_LIBRARY_DIR=/opt/lugrouplib/library
# 1 worker keeps RAM tiny; SQLite + this app do not need more for a small group.
ExecStart=/opt/lugrouplib/venv/bin/gunicorn -w 1 -b 0.0.0.0:8000 app:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
